# handheld/systemd/mirrorctl@.service
[Unit]
Description=mirrorctl status reporter (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/RaspberryPiSystem_001/handheld
Environment=PYTHONPATH=/home/%i/RaspberryPiSystem_001
ExecStart=/bin/bash -c '/home/%i/RaspberryPiSystem_001/handheld/.venv/bin/python /home/%i/RaspberryPiSystem_001/handheld/scripts/mirrorctl.py status'
Restart=on-failure
RestartSec=900

[Install]
WantedBy=multi-user.target
