# scripts/server/toolmaster/systemd/usb-dist-sync@.service
[Unit]
Description=Toolmaster USB dist sync (%I)
Requires=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
Environment=TOOLMASTER_CLIENT_HOME={{CLIENT_HOME}}
Environment=USB_DIST_LABEL=TM-DIST
Environment=RUN_IMPORTER_AFTER_SYNC=1
Environment=IMPORTER_BIN=/usr/local/bin/document-importer.sh
ExecStart=/usr/local/bin/tool-dist-sync.sh --device /dev/%I
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
