# scripts/server/toolmaster/systemd/usb-backup@.service
[Unit]
Description=Toolmaster USB backup export (%I)
Requires=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
Environment=SNAPSHOT_DIR=/srv/RaspberryPiSystem_001/toolmaster/snapshots
Environment=BACKUP_RETENTION=4
Environment=USB_BACKUP_LABEL=TM-BACKUP
ExecStart=/usr/local/bin/tool-backup-export.sh --device /dev/%I
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
